# alien/btrfs-tools-snapshots/etc/cron.daily/btrfs-tools-snapshots
#!/bin/sh

btrfs-tools-snapshots --auto --trigger cron.daily
